# roles/jenkins_ci_backup/templates/jenkins_backup.service.j2
[Unit]
Description=Backup the Jenkins CI service data

[Service]
Type=simple
ExecStart={{ jenkins_home | mandatory }}/backups/jenkins_backup.sh

[Install]
WantedBy=multi-user.target