# roles/backup_client/templates/backup_client.timer.j2
[Unit]
Description=backup client schedule

[Timer]
OnCalendar=*-*-* 04:45:00
#AccuracySec=2 h 30 min
Unit=backup_client.service

[Install]
WantedBy=multi-user.target