# roles/backup_client/templates/backup_client.service.j2
[Unit]
Description=backup client

[Service]
Type=simple
ExecStart=/usr/local/sbin/backup_client

[Install]
WantedBy=multi-user.target